# systemd/agent-forge-polling.service
[Unit]
Description=Agent-Forge Polling Service
After=network.target
Documentation=https://github.com/m0nk111/agent-forge

[Service]
Type=simple
User=flip
Group=flip
WorkingDirectory=/home/flip/agent-forge

# Python environment
Environment="PYTHONPATH=/home/flip/agent-forge"
Environment="PATH=/usr/bin:/usr/local/bin"

# GitHub tokens - Load from environment file
# Create /etc/agent-forge/tokens.env with:
#   BOT_GITHUB_TOKEN=your_token_here
#   CODEAGENT_GITHUB_TOKEN=your_token_here
#   BOT_GITHUB_USERNAME=your_bot_username
#   BOT_GITHUB_EMAIL=your_bot_email
#   CODEAGENT_GITHUB_USERNAME=your_codeagent_username
#   CODEAGENT_GITHUB_EMAIL=your_codeagent_email
#   GITHUB_TOKEN=${BOT_GITHUB_TOKEN}
#   GITHUB_USERNAME=${BOT_GITHUB_USERNAME}
EnvironmentFile=-/etc/agent-forge/tokens.env

# Polling configuration (5 minutes interval)
Environment="POLLING_INTERVAL=300"
Environment="POLLING_REPOS=m0nk111/agent-forge"

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=agent-forge-polling

# Start polling service
ExecStart=/usr/bin/python3 -m engine.runners.polling_service

# Restart policy
Restart=always
RestartSec=30

# Resource limits
MemoryMax=512M
CPUQuota=50%

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/flip/agent-forge

# Timeouts
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
