# systemd/agent-forge.service
[Unit]
Description=Agent-Forge Autonomous Development Platform
Documentation=https://github.com/m0nk111/agent-forge
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
User=agent-forge
Group=agent-forge
WorkingDirectory=/opt/agent-forge

# Environment
Environment="PYTHONUNBUFFERED=1"

# Python environment and command
ExecStart=/opt/agent-forge/scripts/start-service.sh
ExecReload=/bin/kill -HUP $MAINPID

# Restart policy
Restart=on-failure
RestartSec=10
StartLimitBurst=5

# Resource limits
CPUQuota=80%
MemoryMax=4G
TasksMax=100

# Watchdog
WatchdogSec=60
NotifyAccess=main

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/agent-forge/logs /opt/agent-forge/data /opt/agent-forge/config
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictNamespaces=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=agent-forge

# Graceful shutdown
TimeoutStopSec=30
KillMode=mixed
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
