# systemd/agent-forge-polling.service
[Unit]
Description=Agent-Forge Polling Service
After=network.target
Documentation=https://github.com/m0nk111/agent-forge

[Service]
Type=simple
User=flip
Group=flip
WorkingDirectory=/home/flip/agent-forge

# Python environment
Environment="PYTHONPATH=/home/flip/agent-forge"
Environment="PATH=/usr/bin:/usr/local/bin"

# GitHub tokens (configure these!)
# Option 1: Bot-specific token (recommended for automated operations)
Environment="BOT_GITHUB_TOKEN=your_bot_token_here"

# Option 2: Standard GitHub token (fallback)
# Environment="GITHUB_TOKEN=your_token_here"

# Polling configuration
Environment="POLLING_INTERVAL=300"  # 5 minutes
Environment="POLLING_REPOS=m0nk111/agent-forge"
Environment="GITHUB_USERNAME=m0nk111-bot"  # Or your bot account name

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=agent-forge-polling

# Start polling service
ExecStart=/usr/bin/python3 -m engine.runners.polling_service

# Restart policy
Restart=always
RestartSec=30

# Resource limits
MemoryMax=512M
CPUQuota=50%

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/flip/agent-forge

# Timeouts
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
